# nova-bridge.service
[Unit]
Description=Nova WebSocket Bridge for OpenClaw
After=network.target openclaw-gateway.service

[Service]
Type=simple
EnvironmentFile=/home/ec2-user/nova-bridge/.env
ExecStart=/home/ec2-user/.local/bin/nova-bridge
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
